# root/etc/services.d/ICEcoder/run
#!/usr/bin/with-contenv bash

cd /app/ICEcoder || exit

exec \
	s6-setuidgid abc \
	php -S 0.0.0.0:8080
